# spaztick.service
# Systemd unit for Spaztick: Telegram <-> Ollama bridge + config web UI
# Start at boot (headless Ubuntu): enable the service after editing paths below.
#
# 1. Edit SPAZTICK_ROOT below (or sed when copying), and User= if not your user.
# 2. sudo cp spaztick.service /etc/systemd/system/
# 3. sudo systemctl daemon-reload && sudo systemctl enable spaztick && sudo systemctl start spaztick
# 4. Logs: journalctl -u spaztick -f

[Unit]
Description=Spaztick Telegram-Ollama bridge and config UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace with your install path (e.g. /home/ubuntu/spaztick)
WorkingDirectory=/path/to/spaztick
ExecStart=/path/to/spaztick/.venv/bin/python -m run
# Run as your user (replace 'ubuntu' with your username, or remove to run as root)
User=ubuntu
Restart=on-failure
RestartSec=5

[Install]
# Start at boot (multi-user = normal boot on Ubuntu)
WantedBy=multi-user.target
